# {{cookiecutter.project_slug}}/config/supervisor.conf
[program:uwsgi]
command = /usr/bin/make --directory /app start

[program:nginx]
command = /usr/sbin/nginx
