# {{cookiecutter.project_slug}}/config/supervisor.conf
[program:app-uwsgi]
command = /usr/bin/make --directory /app start

[program:nginx-app]
command = /usr/sbin/nginx
